# systemd/.config/systemd/user/evolution-source-registry.service.d/override.conf
[Service]
ExecCondition=/bin/sh -c 'echo "$XDG_CURRENT_DESKTOP" | grep -iq "GNOME"'
